# rel/service/leofs-epmd.socket
[Unit]
Description=Erlang Port Mapper Daemon Activation Socket

[Socket]
ListenStream=0.0.0.0:4369
Accept=false

[Install]
WantedBy=sockets.target
